# ikev2/crontab
0 20 * * * /bin/bash /update-iptables.sh -d
