# ss/crontab
0 20 * * 3,7 /bin/bash /update-iptables.sh -d
